# Scripts/send-email.service
[Unit]
Description=Send Daily Email via Gmail SMTP

[Service]
Type=oneshot
ExecStart=/usr/local/bin/send_email.sh
